# .claude/dispatcher/launchd/com.cklinker.emf-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cklinker.emf-status</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/craigklinker/GitHub/emf/.claude/dispatcher/status.sh > /tmp/emf-status.txt 2>&1</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>EMF_QUEUE_REPO</key>
    <string>/Users/craigklinker/GitHub/emf-queue</string>
    <key>WORKER_HOST</key>
    <string>worker-01</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/emf-status.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/emf-status.launchd.log</string>
</dict>
</plist>
