# .claude/dispatcher/launchd/com.cklinker.emf-planner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cklinker.emf-planner</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/craigklinker/GitHub/emf/.claude/planner/planner.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>EMF_REPO</key>
    <string>/Users/craigklinker/GitHub/emf</string>
    <key>EMF_QUEUE_REPO</key>
    <string>/Users/craigklinker/GitHub/emf-queue</string>
    <key>PATH</key>
    <string>/Users/craigklinker/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/emf-planner.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/emf-planner.launchd.log</string>
</dict>
</plist>
