# .claude/dispatcher/observability/promtail/promtail.service
[Unit]
Description=Promtail (ships emf-dispatcher logs to homelab Loki)
Documentation=https://github.com/cklinker/emf/blob/main/.claude/dispatcher/observability/README.md
After=network-online.target emf-dispatcher.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/promtail -config.file=/etc/promtail/config.yml
Restart=always
RestartSec=10
TimeoutStartSec=30

# Promtail needs to read /var/log/emf-dispatcher/* (group-readable as
# craig:craig per the dispatcher unit) and journald (`journalctl` group
# membership).
SupplementaryGroups=craig systemd-journal

StandardOutput=journal
StandardError=journal
SyslogIdentifier=promtail

[Install]
WantedBy=multi-user.target
